# program_supervisord.conf
[program:foo]
directory=/home/user/projects/project
command=/home/user/.virtualenvs/project/bin/python foo.py
redirect_stderr=true
stdout_logfile=/path/to/foo.log
